# docker/etc/supervisor/services/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm -R -F
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
